# system/etc/profile.d/append-termux-path.sh
PREFIX=/data/data/com.termux/files/usr
append_path $PREFIX/bin ||
unset PREFIX
